# files/etc/uci-defaults/99-default-root-password
#!/bin/sh

if command -v chpasswd >/dev/null 2>&1; then
	echo 'root:password' | chpasswd
else
	passwd root <<'EOF'
password
password
EOF
fi

exit 0
